# crontab
*/20 * * * * python3 /feedme.py /searches.txt /srv/http/index.xml
