# crontab
*/20 * * * * /venv/bin/python /feedme.py /searches.txt /srv/http/index.xml
