# crontab
6 6 * * * /venv/bin/python /feedme.py /searches.txt /searches.pickle /srv/http/index.xml
